# supervisor.conf
[supervisord]
nodaemon=true

[program:bot]
command=npm run run
directory=/app/utilly

[program:server]
command=/app/server/main -base-web-url=/app/web
